# templates/default/init.elasticsearch.conf.erb
# ElasticSearch
#
# This file is based on https://gist.github.com/rbscott/1052015
# It is developed for upstart v0.6.5 (an old version) because CentOS ships with an old upstart version by default.

description     "ElasticSearch"

start on (net-device-up
and local-filesystems
and runlevel [2345])

stop on runlevel [016]

respawn limit 10 5

env ES_HOME="<%= node['elasticsearch']['dir'] %>"
env ES_MIN_MEM=<%= node['elasticsearch']['min_mem'] %>
env ES_MAX_MEM=<%= node['elasticsearch']['max_mem'] %>
env DAEMON="<%= node['elasticsearch']['dir'] %>/bin/elasticsearch"
env DATA_DIR="<%= node['elasticsearch']['data_dir'] %>"
env CONFIG_DIR="<%= node['elasticsearch']['conf_dir'] %>"
env LOG_DIR="<%= node['elasticsearch']['logs_dir'] %>"
env WORK_DIR="<%= node['elasticsearch']['work_dir'] %>"

console output

script

# Override our defaults with user defaults:
[ -f /etc/default/elasticsearch ] && . /etc/default/elasticsearch

sudo -u <%= node['elasticsearch']['user'] %> -s /bin/sh <%= node['elasticsearch']['dir'] %>/bin/elasticsearch -f -Des.path.conf=$CONFIG_DIR -Des.path.home=$ES_HOME -Des.path.logs=$LOG_DIR -Des.path.data=$DATA_DIR -Des.path.work=$WORK_DIR
end script